# systemd/rigforge-api.socket.template
[Unit]
Description=RigForge sister API socket (read-only stats superset, #99)

[Socket]
ListenStream=$API_BIND:$API_PORT
Accept=yes
MaxConnections=16
MaxConnectionsPerSource=8

[Install]
WantedBy=sockets.target
