# systemd/rigforge-autotune.timer.template
[Unit]
Description=Periodic RigForge autotune

[Timer]
OnCalendar=$AUTOTUNE_ONCALENDAR
Persistent=true

[Install]
WantedBy=timers.target
